# packaging/com.birdframe.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key><string>com.birdframe</string>
  <key>ProgramArguments</key>
  <array>
    <string>__UV__</string>
    <string>run</string>
    <string>birdframe</string>
  </array>
  <key>WorkingDirectory</key><string>__WORKDIR__</string>
  <key>RunAtLoad</key><true/>
  <key>KeepAlive</key><true/>
  <key>StandardOutPath</key><string>__HOME__/Library/Logs/birdframe.log</string>
  <key>StandardErrorPath</key><string>__HOME__/Library/Logs/birdframe.log</string>
</dict>
</plist>
